# dashboard/dashboard.service
# Install on the touchscreen Pi:
#   1. Run `whoami` and `pwd` (from inside the dashboard/ folder) first.
#   2. Replace REPLACE_WITH_YOUR_USERNAME below with that username on BOTH
#      the User= line and the WorkingDirectory= path (Raspberry Pi OS
#      Bookworm no longer creates a default "pi" user - it's whatever
#      username you picked in Raspberry Pi Imager).
#   3. sudo cp dashboard.service /etc/systemd/system/dashboard.service
#      sudo systemctl daemon-reload
#      sudo systemctl enable --now dashboard.service

[Unit]
Description=Raspy Lab kiosk dashboard server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_WITH_YOUR_USERNAME
WorkingDirectory=/home/REPLACE_WITH_YOUR_USERNAME/raspy-lab/dashboard
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=3
Environment=PORT=8080

[Install]
WantedBy=multi-user.target
